# wordpress-docker-nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:php7-fpm]
command=/usr/sbin/php-fpm7.0 --nodaemonize
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" 
autostart=true
autorestart=true

